# dropbox_backup/rootfs/etc/s6-overlay/s6-rc.d/dropbox-backup/run
#!/command/with-contenv bashio

exec python3 /app/run.py
